# cookbooks/zeit-app-cache/templates/varnishncsa.service.erb
[Service]
ExecStart=
ExecStart=/usr/bin/varnishncsa -D -a -w /var/log/varnish/varnishncsa.log -F '<%= node['varnish']['ncsa_format'].gsub('%', '%%')%>'
